# deploy/modelgate.service
# Place into /etc/systemd/system/modelgate.service (or your SERVICE_NAME)
[Unit]
Description=modelgate
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=modelgate
Group=modelgate
WorkingDirectory=/opt/modelgate/current

# Optional env overrides (e.g. UPSTREAM_BASE_URL); upstream API key lives in shared/config.toml
EnvironmentFile=-/etc/modelgate/modelgate.env

# Optional: override tracing log dir (same as config [logging].tracing_log_dir)
# Environment=TRACING_LOG_DIR=/opt/modelgate/shared/logs
# Optional: tracing level (default info if unset in code)
# Environment=RUST_LOG=info

ExecStart=/opt/modelgate/current/modelgate
Restart=always
RestartSec=2

# Hardening (tweak as needed)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/modelgate

[Install]
WantedBy=multi-user.target

